# deploy/supervisor.conf
; Winning Products - Supervisor queue workers
; Place in /etc/supervisor/conf.d/ or include from main config

[program:winning-connectors]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/winning-products/artisan queue:work redis --queue=connectors --tries=3 --max-time=3600
autostart=true
autorestart=true
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/winning-connectors.log

[program:winning-ingestion]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/winning-products/artisan queue:work redis --queue=ingestion --tries=2 --max-time=3600
autostart=true
autorestart=true
numprocs=4
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/winning-ingestion.log

[program:winning-scoring]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/winning-products/artisan queue:work redis --queue=scoring --tries=2 --max-time=3600
autostart=true
autorestart=true
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/winning-scoring.log

[program:winning-leaderboard]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/winning-products/artisan queue:work redis --queue=leaderboard --tries=2 --max-time=3600
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/winning-leaderboard.log

[program:winning-notifications]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/winning-products/artisan queue:work redis --queue=notifications --tries=2 --max-time=3600
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/winning-notifications.log
